# etc/systemd/system/logmonitor.service
[Unit]
Description=Log Monitoring Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /var/www/html/log_monitor.sh
Environment="PATH=/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin"
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
